# package/sing-box/files/sing-box.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2024 OpenWrt
# sing-box init script

START=99
STOP=10
USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command /usr/bin/sing-box run -c /etc/sing-box/config.json
	procd_set_param respawn
	procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger "sing-box"
}
